# install_modules/fakehwclock/fake-hwclock-tick.service
[Unit]
Description=Save hwclock by timer

[Service]
Type=oneshot
ExecStart=/usr/bin/fake-hwclock.sh tick > /dev/null

