# install_modules/autostart/autostart.service.template
[Unit]
Description=Autostart Oradio
After=multi-user.target

[Service]
# user and script path will be replaced by install script ../autostart.sh
ExecStart=/home/USER/.venv/bin/python3 SCRIPT_PATH/oradio_control.py &
WorkingDirectory=SCRIPT_PATH
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
