# install_resources/hw_version.service.template
# Placeholders will be modified by ../oradio_install.sh script

[Unit]
Description=Generate HW Serial Number for Oradio
After=multi-user.target

[Service]
# Placeholders will be modified by ../hw_version.sh script
Type=oneshot
WorkingDirectory=PLACEHOLDER_PYTHON_PATH
ExecStart=/home/PLACEHOLDER_USER/.venv/bin/python3 PLACEHOLDER_PYTHON_PATH/hw_serial_number_gen.py

[Install]
WantedBy=multi-user.target
