# install_resources/usb-prepare.service
# https://manpages.debian.org/testing/systemd/systemd.service.5.en.html
# Placeholders will be modified by ../oradio_install.sh script

[Unit]
Description=Early task to delete monitor /media/usb_ready
# Make sure this runs before udev is launched
Before=systemd-udevd.service
DefaultDependencies=no
After=sysinit.target

[Service]
Type=oneshot
ExecStart=/bin/rm -f /media/usb_ready
RemainAfterExit=yes

[Install]
# Runs early in the boot process
WantedBy=sysinit.target
